# DOT-config/autostart/homeassistant.desktop
[Desktop Entry]
Type=Application
Exec=/home/gregj/restart-browser.sh